# sys-apps/upstart/files/notify-cgmanager.conf
start on startup
description "Notify cgroup manager address"
task
script
	initctl notify-cgroup-manager-address "unix:path=/sys/fs/cgroup/cgmanager/sock" || true
end script
