# sys-power/cpufrequtils/files/cpufrequtils-init.d-007
#!/sbin/runscript
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

change() {
	ebegin "Running cpufreq-set ${*}"
	/usr/libexec/cpufrequtils-change.sh "${@}" > /dev/null
	eend $?
}

start() {
	change ${START_OPTS} -- ${SYSFS_EXTRA}
}

stop() {
	change "${STOP_OPTS}"
}
